# scripts/io.github.zirconeey.email-summary.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>io.github.zirconeey.email-summary</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/zhourui/Desktop/zirconeey.github.io/scripts/email_summary.sh</string>
    </array>

    <!-- 每天本地时间 09:00 和 21:00 各触发一次；macOS 睡眠错过的会在唤醒后补跑。
         脚本内按当前小时判定 morning / evening 时段，并各自记 lastrun 防重跑。 -->
    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Hour</key><integer>9</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key><integer>21</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
    </array>

    <!-- 登录 / 开机加载时也触发一次，配合脚本的时段+lastrun 判断做「开机补跑」 -->
    <key>RunAtLoad</key>
    <true/>

    <!-- 不长驻；跑完即退 -->
    <key>KeepAlive</key>
    <false/>

    <key>WorkingDirectory</key>
    <string>/Users/zhourui/Desktop/zirconeey.github.io</string>

    <key>StandardOutPath</key>
    <string>/Users/zhourui/Library/Logs/zirconeey-email-summary.launchd.out</string>
    <key>StandardErrorPath</key>
    <string>/Users/zhourui/Library/Logs/zirconeey-email-summary.launchd.err</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/zhourui</string>
        <!-- 代理由 scripts/email_summary.sh 探测后动态决定（国内套代理 / 美国直连） -->
    </dict>
</dict>
</plist>
